# systemd/user/spotifyd.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.config/spotifyd/run-spotifyd.sh --no-daemon
