# systemd/user/dropbox.service.d/override.conf
[Service]
ExecStart=
ExecStart=%h/.dropbox-dist/dropboxd
