# systemd/user/gpg-agent-ssh.socket.d/override.conf
[Socket]
ExecStartPost=/usr/bin/systemctl --user set-environment SSH_AUTH_SOCK=%t/gnupg/S.gpg-agent.ssh