# systemd/user/emacs.service.d/override.conf
[Service]
Type=notify

# Force the use of Xwayland for now
Environment=GDK_BACKEND=x11