# systemd/user/applications.target
[Unit]
Description=Start user applications
IgnoreOnIsolate=true
Requires=dbus.socket xorg.target window-manager.service
