# systemd/user/awesome.service
[Unit]
Description=Awesome WM
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/awesome
#Restart=on-failure

[Install]
Alias=window-manager.service
