# systemd/user/awesome.service
[Unit]
Description=Awesome WM

[Service]
ExecStart=/usr/bin/awesome
Restart=always

[Install]
WantedBy=default.target
