# systemd/user/compton.service
[Unit]
Description=X compositor
Requires=xorg.target

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=xorg.target
