# systemd/user/environment.target
[Unit]
Description=Set session environment variables
IgnoreOnIsolate=true
Requires=dbus.socket
RequiredBy=xorg.target
Before=xorg.target
