# systemd/user/mpd.service
[Unit]
Description=Music Player Daemon

[Service]
ExecStart=/usr/bin/mpd --no-daemon %h/.config/mpd.conf

[Install]
WantedBy=default.target
