# systemd/user/mpris2-info.service
[Unit]
Description=MPRIS2 notifications service
Requires=dbus.socket xorg.target

[Service]
BusName=net.schnouki.Mpris2Info
ExecStart=/home/schnouki/bin/mpris2-info

[Install]
WantedBy=xorg.target
