# systemd/user/mpris2-info.service
[Unit]
Description=MPRIS2 notifications service
Requires=graphical-session.target

[Service]
Type=dbus
BusName=net.schnouki.Mpris2Info
ExecStart=/home/schnouki/bin/mpris2-info --journal
Restart=on-failure

[Install]
WantedBy=applications.target
