# systemd/user/mullvad-vpn.service
[Unit]
Description=Mullvad VPN
Requires=graphical-session.target

[Service]
ExecStart=/usr/bin/mullvad-vpn

[Install]
WantedBy=applications.target
