# systemd/user/picom.service
[Unit]
Description=X compositor
Requires=xorg.target

[Service]
ExecStart=/usr/bin/picom

[Install]
WantedBy=applications.target
