# systemd/user/synclient.service
[Service]
Type=oneshot
ExecStart=/usr/bin/synclient TouchpadOff=1

[Install]
WantedBy=default.target
