# systemd/user/wmname.service
[Service]
Type=oneshot
ExecStart=/usr/bin/wmname LG3D

[Install]
WantedBy=default.target
