# systemd/user/wpa_gui.service
[Unit]
Description=WPA Graphical User Interface
Requires=xorg.target

[Service]
ExecStart=/usr/bin/wpa_gui -t

[Install]
WantedBy=xorg.target
