# systemd/user/wpa_gui.service
[Unit]
Description=WPA Graphical User Interface

[Service]
ExecStart=/usr/bin/wpa_gui -t

[Install]
WantedBy=default.target
