# systemd/user/xhost-local-root.service
[Unit]
Description=Allow root to use this display

[Service]
Type=oneshot
ExecStart=/usr/bin/xhost +local:root

[Install]
WantedBy=default.target
