# systemd/user/xorg.target
[Unit]
Description=The Xorg server
Requires=window-manager.service
