# systemd/user/xscreensaver.service
[Unit]
Description=X Screensaver

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash

[Install]
WantedBy=default.target
