# systemd/user/xscreensaver.service
[Unit]
Description=X Screensaver
Requires=xorg.target

[Service]
ExecStart=/usr/bin/xscreensaver -no-splash

[Install]
WantedBy=xorg.target
