# systemd/user/pueued.service.d/override.conf
[Service]
PIDFile=/home/schnouki/.local/share/pueue/pueue.pid
Restart=on-failure